# Definitions/main.mail_server/docker/postfix/1.0.0/proto.sendmail.service
[Unit]
Description=Postfix service running on Docker
After=docker.service
Requires=docker.service

[Service]
Restart=always
ExecStart={{DOCKER.HOME}}/Definitions/Postfix/start.sh

[Install]
WantedBy=multi-user.target